# trusted_computing/distributions/gentoo/keyrona
#!/sbin/runscript

depend() { 
    after tcsd
}

start() {
	ebegin "Starting Keyrona's keyproviderd"
	killall -q keyrona_keyproviderd
	start-stop-daemon --start --background --exec /usr/sbin/keyrona_keyproviderd -- --prune
	eend $?
}

stop() {
	ebegin "Stopping Keyrona's keyproviderd"
	killall -q keyrona_keyproviderd 
	eend 0
}
