# trusted_computing/distributions/gentoo/tpmcrypt
#!/sbin/runscript

depend() { 
    after tcsd
}

start() {
	ebegin "Starting TpmCrypt's keyproviderd"
	killall -q tpmcrypt_keyproviderd
	start-stop-daemon --start --background --exec /usr/sbin/tpmcrypt_keyproviderd -- --prune
	eend $?
}

stop() {
	ebegin "Stopping TpmCrypt's keyproviderd"
	killall -q tpmcrypt_keyproviderd 
	eend 0
}
