# tools/flip/config/templates/etc/rc.local
#!/bin/sh
ifconfig sis0 down
. $NETSTART
